# files/systemd/instance_switch@.service
[Service]
ExecStart=/opt/galera/scripts/instance_switch.sh
User=nagios
StandardInput=socket

[Install]
WantedBy=multi-user.target
