# files/systemd/instance_switch.socket
[Socket]
ListenStream=1234
Accept=true

[Install]
WantedBy=sockets.target
